# rootfs/etc/systemd/system/usbguard.service.d/override.conf
[Service]
ExecStartPre=-mkdir -p /var/log/usbguard
ExecStartPre=-/bin/sh -C -c 'usbguard generate-policy > /etc/usbguard/rules.d/custom.conf'
